# debian-pkg/usr/lib/udev/rules.d/90-tinypilot-hid.rules
SUBSYSTEM=="hidg", ACTION=="add", MODE="0660", GROUP="tinypilot", RUN+="/usr/libexec/tinypilot-privileged/label-hid-dev %M %m %k"
