# debian-pkg/usr/lib/systemd/system/janus.service.d/20-tinypilot.conf
[Unit]
# Give up if we restart on failure 20 times within 5 minutes (300 seconds).
StartLimitIntervalSec=300
StartLimitBurst=20

[Service]
RestartSec=1
